# deploy_tools/mgh.conf
[program:mgh]
command=/home/nanoafrika/mgh_supervisior
user=nanoafrika
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/nanoafrika/logs/mgh.log